# server/image/ansible/playbooks/files/profile.d/env.sh
#
# this file is provisioned by Packer
#
LC_COLLATE=C
LESS="-FMRX"
export LC_COLLATE LESS
