# server/image/ansible/playbooks/files/profile.d/options.sh
#
# this file is provisioned by Packer
#
set -o vi
